# conf/supervisor.conf
[program:mattermost]
directory=/var/www/mattermost
command=/var/www/mattermost/bin/platform
user=www
startsecs=5
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/mattermost.log
