# conf/supervisor.conf
[program:mattermost]
directory=/var/www/mattermost
command=/var/www/mattermost/bin/platform
startsecs=5
autorestart=true
